# supervisor/conf.d/optional-service-template.conf
; supervisor/conf.d/optional-service-template.conf – Template for adding services.
;
; Copy and customise this file to run additional daemons alongside piclaw.
; Keep optional services isolated from piclaw. If the binary is missing, exit 0.

[program:optional-service]
command=/bin/bash -lc 'command -v /usr/bin/example-service >/dev/null || exit 0; exec /usr/bin/example-service --flag value'
user=agent
autostart=false
autorestart=false
startsecs=1
startretries=0
stopsignal=INT
stopwaitsecs=10
stdout_logfile=/var/log/piclaw/optional-service.stdout.log
stderr_logfile=/var/log/piclaw/optional-service.stderr.log
